# crontab
0 3 * * * /home/rbitel/Configuration-Scripts/backup_data.sh > /dev/null
